# packaging/debian/postrm
#!/bin/bash
set -e

case "$1" in
    remove)
        echo ""
        echo "Hefesto - Dualsense4Unix removido."
        echo ""
        echo "Configuracoes do usuario foram preservadas em:"
        echo "  ~/.config/hefesto-dualsense4unix/  (perfis, sessao, preferencias)"
        echo "  ~/.config/hefesto/                 (layout legado pre-rename, se existir)"
        echo "Para remover tambem as configuracoes:"
        echo "  rm -rf ~/.config/hefesto-dualsense4unix/ ~/.config/hefesto/"
        echo ""
        echo "Se o servico systemd ainda estiver habilitado para algum usuario:"
        echo "  systemctl --user disable --now hefesto-dualsense4unix.service"
        ;;
    purge)
        # Limpeza do venv bundlado (apt remove apenas arquivos rastreados;
        # o venv tem .pyc gerados em runtime que ficam orfaos sem isso).
        rm -rf /opt/hefesto-dualsense4unix 2>/dev/null || true
        ;;
    upgrade|failed-upgrade|abort-install|abort-upgrade|disappear) ;;
esac

exit 0
