# packaging/debian/postinst
#!/bin/bash
set -e

case "$1" in
    configure)
        # Reload das regras udev para DualSense
        if command -v udevadm >/dev/null 2>&1; then
            udevadm control --reload-rules
            udevadm trigger --action=change --subsystem-match=usb
        fi

        # Atualiza cache de icones GTK
        if command -v gtk-update-icon-cache >/dev/null 2>&1; then
            gtk-update-icon-cache -q -f /usr/share/icons/hicolor || true
        fi

        # Atualiza base de dados de .desktop
        if command -v update-desktop-database >/dev/null 2>&1; then
            update-desktop-database -q /usr/share/applications || true
        fi

        # NÃO habilita o servico systemd automaticamente — opt-in.
        echo ""
        echo "Hefesto - Dualsense4Unix instalado com sucesso."
        echo ""
        echo "Todas as dependencias Python estao bundladas em"
        echo "  /opt/hefesto-dualsense4unix/venv/"
        echo "(nao requer pip install manual — funciona em qualquer Ubuntu/Pop! >= 22.04)."
        echo ""
        echo "Para habilitar o inicio automatico do daemon:"
        echo "  systemctl --user enable --now hefesto-dualsense4unix.service"
        echo ""
        echo "Broker root hide-hidraw (BROKER-01 — esconde o controle FISICO do"
        echo "jogo, cura de raiz do duplicado; o postinst roda sem sessao e NAO"
        echo "pode renderizar o uid — ative da SUA sessao de usuario) e os modulos"
        echo "DKMS patchados hid-nintendo (Onda T — cura de raiz do probe BT dos"
        echo "controles Nintendo/8BitDo) e rtw88_usb (Onda W — cura de raiz do"
        echo "fantasma USB do dongle WiFi; ambos requerem os pacotes dkms + headers):"
        echo "  sudo /usr/share/hefesto-dualsense4unix/scripts/install-host-udev.sh"
        echo ""
        echo "Ou abra 'Hefesto' no menu de aplicativos."
        ;;
    abort-upgrade|abort-remove|abort-deconfigure) ;;
esac

exit 0
