# assets/systemd/hefesto-bt-health-watchdog.service
[Unit]
# Watchdog do estado "vivo mas doente" do bluetoothd + vigia de bond
# temporário (sprint 2026-07-21 BlueZ, camada 2). Complementa (não substitui)
# Restart=on-failure e WatchdogSec do drop-in 10-hefesto-resilience.conf:
# aqueles cobrem crash e hang; este cobre o daemon SÃO para o systemd que
# recusa devices pareados em loop pós-crash, e o bond Paired-sem-Bonded que
# evapora no disconnect (medido 21-22/07).
Description=Hefesto — watchdog de saúde do bluetoothd e persistência de bonds
After=bluetooth.service

[Service]
Type=oneshot
ExecStart=/usr/local/lib/hefesto-dualsense4unix/bt_health_watchdog.sh
NoNewPrivileges=yes
ProtectHome=yes
PrivateTmp=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
SystemCallArchitectures=native
