# assets/77-dualsense-leds.rules
# Hefesto - Dualsense4Unix — acesso de escrita aos LEDs do DualSense via sysfs.
# FEAT-DSX-LIGHTBAR-SYSFS-01.
#
# O QUE FAZ: torna graváveis (pelo usuário do desktop) os nós de LED que o driver
# de kernel `hid_playstation` expõe para o DualSense — a lightbar RGB
# (`*:rgb:indicator`, atributo multi_intensity) e os 5 LEDs de player
# (`*:white:player-N`, atributo brightness).
#
# POR QUE: escrever nesses nós DELEGA ao kernel a montagem do output report, que
# difere entre USB e Bluetooth (no BT exige seq_tag monotônico + CRC-32). Assim a
# COR da lightbar acende IGUAL em USB e BT. A escrita crua por hidraw (pydualsense)
# usa seq_tag fixo e disputa esses LEDs com o próprio kernel, fazendo a cor "não
# colar" no BT (BUG-MULTI-CONTROLLER-BT-CRC-CONTENTION-01 / lightbar-bt).
#
# Sem esta regra, os nós ficam root:root e o daemon (que roda como usuário) cai no
# caminho pydualsense — sem regressão, mas a cor por BT volta a falhar. Por isso a
# regra e o código sysfs entram JUNTOS.
#
# MECANISMO: RUN+= chmod 0666 garante gravabilidade dos nós de LED (LED é
# inócuo; world-writable aqui não é risco). Reaplicado a cada replug/reconexão
# (ACTION=="add").
#
# ONDA-R (2026-07-19): TAG+="uaccess" foi REMOVIDA desta regra — este arquivo
# é instalado como 77-dualsense-leds.rules, e systemd-logind só converte a TAG
# "uaccess" em ACL para regras numeradas < 73 (é o 73-seat-late.rules, do
# próprio systemd, quem faz essa conversão — ver reference_udev_uaccess_ordem_73
# na memória do projeto). Numa regra ≥73 a TAG fica presente no device mas
# nunca vira ACL: era morta/inerte desde sempre nesta máquina, sem efeito
# funcional (o chmod 0666 abaixo já cobre 100% do acesso). Removida para não
# sugerir uma garantia que a ordem de carga do udev não entrega.

ACTION=="add", SUBSYSTEM=="leds", KERNEL=="*:rgb:indicator", RUN+="/bin/chmod 0666 /sys/class/leds/%k/multi_intensity /sys/class/leds/%k/brightness"
ACTION=="add", SUBSYSTEM=="leds", KERNEL=="*:white:player-*", RUN+="/bin/chmod 0666 /sys/class/leds/%k/brightness"
