# assets/hefesto-dualsense4unix.service
[Unit]
Description=Hefesto - Dualsense4Unix DualSense adaptive trigger daemon
# BUG-DEB-AUTOSTART-WANTEDBY-DEFAULT-01: o daemon usa /dev/hidraw e evdev,
# não depende de DISPLAY/Wayland/X11. Ancorar em default.target garante
# autostart resiliente em qualquer DE (GNOME/KDE/Sway/i3) e elimina race
# entre login e graphical-session activation. After= preservado por ordem.
After=graphical-session.target default.target
# BUG-MULTI-INSTANCE-01: teto anti-loop em Unit. Máximo 3 tentativas em 30s —
# se estourar, systemd marca failed e pára de respawnar.
StartLimitIntervalSec=30
StartLimitBurst=3

[Service]
Type=simple
# AUTOSWITCH-FLOOD-FIX-01: importa o env gráfico (WAYLAND_DISPLAY/DISPLAY) do
# gerenciador de usuário para dentro do serviço, para o autoswitch (perfil por
# app ativo) sair do NullBackend em COSMIC/Wayland — sem isso o daemon nasce
# sem display, o autoswitch fica morto e o window_detect floodava o journal.
# Prefixo '-' = não-fatal se o compositor ainda não exportou nada (o subsystem
# tem fallback via `systemctl --user show-environment` no _ensure_display_env).
ExecStartPre=-/usr/bin/systemctl --user import-environment WAYLAND_DISPLAY DISPLAY
ExecStart=%h/.local/bin/hefesto-dualsense4unix daemon start --foreground
Restart=on-failure
RestartSec=2
# BUG-MULTI-INSTANCE-01: SIGTERM (143) é desligamento limpo — o takeover da
# próxima instância e o "Sair" do tray não devem disparar respawn.
SuccessExitStatus=143 SIGTERM
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
